# com.eventhub.ig-post.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eventhub.ig-post</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/noa/event-hub-sync/event-hub-pipeline/run-ig-post.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/noa/event-hub-sync/event-hub-pipeline/logs/launchd-ig-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noa/event-hub-sync/event-hub-pipeline/logs/launchd-ig-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/noa/event-hub-sync/event-hub-pipeline</string>
</dict>
</plist>
